# R30B1/files/etc/crontabs/root
0 3 * * 2-7 /etc/init.d/zerotier restart
0 3 * * 1 reboot
